# services/js-content-fetcher/js-content-fetcher.service
[Unit]
Description=JavaScript Content Fetcher Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/js-content-fetcher
ExecStart=/usr/bin/node /opt/js-content-fetcher/server.js
Restart=always
RestartSec=5
StandardOutput=append:/var/log/js-content-fetcher.log
StandardError=append:/var/log/js-content-fetcher-error.log

# Security: By default binds to 127.0.0.1 (localhost only)
# Set BIND_ADDRESS=0.0.0.0 only if needed (not recommended for production)
# Environment=BIND_ADDRESS=127.0.0.1

[Install]
WantedBy=multi-user.target
